# deploy/fitlit.service
[Unit]
Description=FitLit Google Health fetcher and local dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__FITLIT_USER__
WorkingDirectory=__FITLIT_ROOT__
EnvironmentFile=-__FITLIT_ROOT__/.env
Environment="PATH=__FITLIT_PATH__"
ExecStart=__UV_PATH__ run uvicorn fitlit.server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
